# ops/systemd/gaokao-backup.timer
[Unit]
Description=Run gaokao backup snapshot daily

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
Unit=gaokao-backup.service

[Install]
WantedBy=timers.target
